# openrc/shadowdhcp
#!/sbin/openrc-run
# OpenRC init script for shadowdhcp DHCP server

supervisor=supervise-daemon
name="shadowdhcp"
description="Reservation-only DHCPv4/DHCPv6 server"

command="/usr/local/bin/shadowdhcp"
command_args="--configdir /etc/shadowdhcp"
command_user="shadowdhcp:shadowdhcp"

capabilities="^cap_net_bind_service"

output_log="/var/log/shadowdhcp/shadowdhcp.log"
error_log="/var/log/shadowdhcp/shadowdhcp.log"

extra_started_commands="reload"
description_reload="Reload DHCP reservations"

depend() {
    need net
    after firewall
}

start_pre() {
    checkpath --directory --owner shadowdhcp:shadowdhcp --mode 0750 /var/log/shadowdhcp
}

reload() {
    ebegin "Reloading ${name}"
    ${supervisor} ${RC_SVCNAME} --signal HUP
    eend $?
}